# polaris-controller/src/main/dist/bin/polaris.service
[Unit]
Description=The RACTF challenge scheduler
Requires=consul.service
After=consul.service

[Service]
WorkingDirectory=/opt/polaris
ExecStart=/bin/java -jar /opt/polaris/polaris.jar server /etc/polaris/polaris.yml
Type=simple
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target